# src/fts3rest/fts-rest-graceful.cron
12 */6 * * * root /usr/sbin/apachectl graceful >& /dev/null
